# exec/systemd/qemu-cpu@.service
[Unit]
Description="Per CPU Watcher %I"
Slice=qemu 

[Service]
Type=Exec
PIDFile=/var/run/qemu/%I
