# systemd/kvm-local-backup.timer
[Unit]
Description=KVM Backuo

[Timer]
OnCalendar=Mon-Fri *-*-* 07:00:00
Persistent=true

[Install]
WantedBy=basic.target
